# ETM/dataclean/dataclean-api.service
[Unit]
Description=DataClean API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/dataclean
Environment="PATH=/opt/dataclean/venv/bin"
ExecStart=/opt/dataclean/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dataclean-api

[Install]
WantedBy=multi-user.target
